# SpringBoot-ALF/debian/prerm.sh
#!/bin/sh
service springboot-alf stop

#DEBHELPER#

systemctl disable springboot-alf